# airootfs/etc/skel/.config/autostart/InstallSystem.desktop
[Desktop Entry]
Name=Install_RenOS
Exec=sh -c "xhost si:localuser:root && kgx -e doas /home/live/.local/bin/tir"
Type=Application
