# airootfs/etc/skel/.config/autostart/InstallSystem.desktop
[Desktop Entry]
Name=Install_RenOS
Exec=kgx -e ~/.local/bin/tir
Type=Application
